# e2e-pipeline/offline/supervisord.conf
; ============================================================================
; Zero-Touch QA All-in-One — supervisord 프로세스 매니페스트 (Mac 하이브리드)
; 단일 컨테이너 내부에서 10개 프로세스를 관리한다.
; (내부 Ollama 제거 — 호스트 Mac Ollama / 내부 jenkins-agent 제거 — 호스트 Mac agent)
;
; 기동 순서 (priority 낮을수록 먼저):
;   100: 상태 저장소 (postgresql, redis, qdrant)
;   200: dify-plugin-daemon (postgres 준비 후)
;   300: dify api/worker/worker_beat/web
;   400: nginx (api/web 준비 후)
;   500: jenkins (마스터)
;   (jenkins-agent 는 호스트 Mac 에서 직접 JNLP 연결 — 컨테이너 내부 프로그램 없음)
;
; 로그는 /data/logs/<svc>.log 로 분리 저장 (볼륨 지속).
; ============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/data/logs/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ────────────────────────────────────────────────────────────────────────────
; 상태 저장소 (priority 100)
; ────────────────────────────────────────────────────────────────────────────

[program:postgresql]
command=/usr/lib/postgresql/15/bin/postgres -D /data/pg -c config_file=/data/pg/postgresql.conf
user=postgres
priority=100
autostart=true
autorestart=unexpected
stdout_logfile=/data/logs/postgresql.log
stderr_logfile=/data/logs/postgresql.err.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB

[program:redis]
command=/usr/bin/redis-server --dir /data/redis --appendonly yes --bind 127.0.0.1 --port 6379
user=redis
priority=100
autostart=true
autorestart=unexpected
stdout_logfile=/data/logs/redis.log
stderr_logfile=/data/logs/redis.err.log

[program:qdrant]
; qdrant 1.8.3 의 CLI 는 --storage-dir 플래그를 지원하지 않는다.
; 스토리지 경로는 QDRANT__STORAGE__STORAGE_PATH / ..__SNAPSHOTS_PATH 환경변수로 주입.
command=/usr/local/bin/qdrant
environment=
    QDRANT__SERVICE__HTTP_PORT="6333",
    QDRANT__SERVICE__HOST="127.0.0.1",
    QDRANT__STORAGE__STORAGE_PATH="/data/qdrant",
    QDRANT__STORAGE__SNAPSHOTS_PATH="/data/qdrant/snapshots",
    QDRANT__TELEMETRY_DISABLED="true"
priority=100
autostart=true
autorestart=unexpected
stdout_logfile=/data/logs/qdrant.log
stderr_logfile=/data/logs/qdrant.err.log

; [program:ollama] 제거됨 — Mac 브랜치는 호스트 Ollama 를 사용하므로 컨테이너
; 내부 ollama 데몬을 띄우지 않는다. 이 이미지의 기동 시점에 `supervisorctl status`
; 에는 ollama 항목이 아예 나타나지 않으며, `ollama` CLI 또한 컨테이너 안에서
; 호출 불가 (바이너리 미포함).

; ────────────────────────────────────────────────────────────────────────────
; Dify 플러그인 데몬 (priority 200)
; ────────────────────────────────────────────────────────────────────────────

[program:dify-plugin-daemon]
; upstream 이미지의 실제 바이너리는 /app/main (entrypoint.sh 가 `exec /app/main`).
; SERVER_KEY 는 dify-api 의 PLUGIN_DAEMON_KEY 와 동일해야 함 (401 방지).
command=/opt/dify-plugin-daemon/main
directory=/opt/dify-plugin-daemon
environment=
    PATH="/usr/local/bin:/usr/bin:/bin",
    PLATFORM="local",
    SERVER_KEY="QaHbTe77CtuXmsfyhR7+vRjI/+XiGgmVDLZuA+Wf6RTUu8l0mq5kd20MOcc=",
    UV_PATH="/usr/local/bin/uv",
    TIKTOKEN_CACHE_DIR="/opt/dify-plugin-daemon/.tiktoken",
    DB_HOST="127.0.0.1",
    DB_PORT="5432",
    DB_USERNAME="postgres",
    DB_PASSWORD="difyai123456",
    DB_DATABASE="dify_plugin",
    REDIS_HOST="127.0.0.1",
    REDIS_PORT="6379",
    SERVER_PORT="5002",
    PLUGIN_WORKING_PATH="/data/dify/plugins/cwd",
    PLUGIN_STORAGE_TYPE="local",
    PLUGIN_STORAGE_LOCAL_ROOT="/data/dify/plugins",
    PLUGIN_INSTALLED_PATH="plugin",
    PLUGIN_PACKAGE_CACHE_PATH="plugin_packages",
    PLUGIN_MEDIA_CACHE_PATH="assets",
    PLUGIN_REMOTE_INSTALLING_HOST="0.0.0.0",
    PLUGIN_REMOTE_INSTALLING_PORT="5003",
    MAX_PLUGIN_PACKAGE_SIZE="52428800",
    PYTHON_ENV_INIT_TIMEOUT="120",
    PLUGIN_MAX_EXECUTION_TIMEOUT="600",
    FORCE_VERIFYING_SIGNATURE="false",
    DIFY_INNER_API_URL="http://127.0.0.1:5001",
    DIFY_INNER_API_KEY="QaHbTe77CtuXmsfyhR7+vRjI/+XiGgmVDLZuA+Wf6RTUu8l0mq5kd20MOcc=",
    PPROF_ENABLED="false",
    HOST_WEBHOOK="0.0.0.0"
priority=200
autostart=true
autorestart=unexpected
startsecs=10
stdout_logfile=/data/logs/dify-plugin-daemon.log
stderr_logfile=/data/logs/dify-plugin-daemon.err.log

; ────────────────────────────────────────────────────────────────────────────
; Dify 애플리케이션 서비스 (priority 300)
; Dify 환경변수는 docker-compose.yaml 의 x-shared-api-worker-env 앵커와 동일한
; 조합이되 host 값을 전부 127.0.0.1 로 치환한 것.
; ────────────────────────────────────────────────────────────────────────────

[program:dify-api]
; upstream 이미지 레이아웃: /app/api/ (WORKDIR) → allinone 에서 /opt/dify-api/api/
; /app/api 심볼릭 링크 (Dockerfile) 덕분에 기존 upstream ENV (/app/api/.venv 등) 도 작동.
command=/bin/bash /opt/dify-api/api/docker/entrypoint.sh
directory=/opt/dify-api/api
environment=
    PATH="/opt/dify-api/api/.venv/bin:/usr/local/bin:/usr/bin:/bin",
    VIRTUAL_ENV="/opt/dify-api/api/.venv",
    FLASK_APP="app.py",
    TIKTOKEN_CACHE_DIR="/opt/dify-api/api/.tiktoken_cache",
    NLTK_DATA="/usr/local/share/nltk_data",
    LANG="C.UTF-8",
    LC_ALL="C.UTF-8",
    PYTHONIOENCODING="utf-8",
    EDITION="SELF_HOSTED",
    DEPLOY_ENV="PRODUCTION",
    MODE="api",
    CONSOLE_API_URL="",
    CONSOLE_WEB_URL="",
    SERVICE_API_URL="",
    APP_WEB_URL="",
    DB_HOST="127.0.0.1",
    DB_PORT="5432",
    DB_USERNAME="postgres",
    DB_PASSWORD="difyai123456",
    DB_DATABASE="dify",
    REDIS_HOST="127.0.0.1",
    REDIS_PORT="6379",
    REDIS_PASSWORD="",
    CELERY_BROKER_URL="redis://127.0.0.1:6379/1",
    VECTOR_STORE="qdrant",
    QDRANT_URL="http://127.0.0.1:6333",
    STORAGE_TYPE="opendal",
    OPENDAL_SCHEME="fs",
    OPENDAL_FS_ROOT="/data/dify/storage",
    PLUGIN_DAEMON_URL="http://127.0.0.1:5002",
    PLUGIN_DAEMON_KEY="QaHbTe77CtuXmsfyhR7+vRjI/+XiGgmVDLZuA+Wf6RTUu8l0mq5kd20MOcc=",
    INNER_API_KEY_FOR_PLUGIN="QaHbTe77CtuXmsfyhR7+vRjI/+XiGgmVDLZuA+Wf6RTUu8l0mq5kd20MOcc=",
    MIGRATION_ENABLED="true",
    LOG_LEVEL="INFO",
    SECRET_KEY="dify-allinone-placeholder-CHANGE-ME-via-env",
    SERVER_WORKER_AMOUNT="1",
    SERVER_WORKER_CLASS="gevent",
    SERVER_WORKER_CONNECTIONS="1000",
    GUNICORN_TIMEOUT="360"
priority=300
; autostart=false — plugin-daemon race 방지용. entrypoint-allinone.sh 가
; dify-plugin-daemon 의 /health/check 가 200 응답한 뒤 명시적으로
; `supervisorctl start dify-api` 를 호출한다. 과거 cold-start 시 두 서비스가
; 동시에 spawn 되면 plugin-daemon 의 초기 crash → restart 순간에 dify-api
; 의 Flask import 가 plugin-daemon 에 blocking call 을 시도해 "Booting worker"
; 를 못 찍고 영구 대기하는 패턴이 재현됐음.
autostart=false
autorestart=unexpected
startsecs=15
stdout_logfile=/data/logs/dify-api.log
stderr_logfile=/data/logs/dify-api.err.log

[program:dify-worker]
command=/bin/bash /opt/dify-api/api/docker/entrypoint.sh
directory=/opt/dify-api/api
environment=
    PATH="/opt/dify-api/api/.venv/bin:/usr/local/bin:/usr/bin:/bin",
    VIRTUAL_ENV="/opt/dify-api/api/.venv",
    FLASK_APP="app.py",
    TIKTOKEN_CACHE_DIR="/opt/dify-api/api/.tiktoken_cache",
    NLTK_DATA="/usr/local/share/nltk_data",
    LANG="C.UTF-8",
    LC_ALL="C.UTF-8",
    PYTHONIOENCODING="utf-8",
    EDITION="SELF_HOSTED",
    DEPLOY_ENV="PRODUCTION",
    MODE="worker",
    DB_HOST="127.0.0.1",
    DB_PORT="5432",
    DB_USERNAME="postgres",
    DB_PASSWORD="difyai123456",
    DB_DATABASE="dify",
    REDIS_HOST="127.0.0.1",
    REDIS_PORT="6379",
    CELERY_BROKER_URL="redis://127.0.0.1:6379/1",
    VECTOR_STORE="qdrant",
    QDRANT_URL="http://127.0.0.1:6333",
    STORAGE_TYPE="opendal",
    OPENDAL_SCHEME="fs",
    OPENDAL_FS_ROOT="/data/dify/storage",
    PLUGIN_DAEMON_URL="http://127.0.0.1:5002",
    PLUGIN_DAEMON_KEY="QaHbTe77CtuXmsfyhR7+vRjI/+XiGgmVDLZuA+Wf6RTUu8l0mq5kd20MOcc=",
    INNER_API_KEY_FOR_PLUGIN="QaHbTe77CtuXmsfyhR7+vRjI/+XiGgmVDLZuA+Wf6RTUu8l0mq5kd20MOcc=",
    LOG_LEVEL="INFO",
    SECRET_KEY="dify-allinone-placeholder-CHANGE-ME-via-env"
priority=300
autostart=true
autorestart=unexpected
startsecs=15
stdout_logfile=/data/logs/dify-worker.log
stderr_logfile=/data/logs/dify-worker.err.log

[program:dify-worker-beat]
command=/bin/bash /opt/dify-api/api/docker/entrypoint.sh
directory=/opt/dify-api/api
environment=
    PATH="/opt/dify-api/api/.venv/bin:/usr/local/bin:/usr/bin:/bin",
    VIRTUAL_ENV="/opt/dify-api/api/.venv",
    FLASK_APP="app.py",
    TIKTOKEN_CACHE_DIR="/opt/dify-api/api/.tiktoken_cache",
    NLTK_DATA="/usr/local/share/nltk_data",
    LANG="C.UTF-8",
    LC_ALL="C.UTF-8",
    PYTHONIOENCODING="utf-8",
    EDITION="SELF_HOSTED",
    DEPLOY_ENV="PRODUCTION",
    MODE="beat",
    DB_HOST="127.0.0.1",
    DB_PORT="5432",
    DB_USERNAME="postgres",
    DB_PASSWORD="difyai123456",
    DB_DATABASE="dify",
    REDIS_HOST="127.0.0.1",
    REDIS_PORT="6379",
    CELERY_BROKER_URL="redis://127.0.0.1:6379/1",
    LOG_LEVEL="INFO",
    SECRET_KEY="dify-allinone-placeholder-CHANGE-ME-via-env"
priority=300
autostart=true
autorestart=unexpected
startsecs=15
stdout_logfile=/data/logs/dify-worker-beat.log
stderr_logfile=/data/logs/dify-worker-beat.err.log

[program:dify-web]
; upstream 이미지 레이아웃: /app/web/ → allinone 에서 /opt/dify-web/web/ (+/app/web 심볼릭)
; node 는 Dockerfile 에서 nodejs.org glibc 빌드를 /usr/local/bin/node 에 설치.
; MARKETPLACE_API_URL / MARKETPLACE_URL 은 Next.js 16 env validator 가 URL 포맷
; 필수라서 오프라인에서도 더미 URL 을 반드시 넣어야 한다 (미설정 시 /api/v1 로
; 접두사만 남아 invalid_format 에러 → /install 500).
command=/bin/sh -c "cd /opt/dify-web/web && /opt/dify-web/web/entrypoint.sh"
directory=/opt/dify-web/web
environment=
    PATH="/usr/local/bin:/usr/bin:/bin",
    NODE_ENV="production",
    HOSTNAME="127.0.0.1",
    PORT="3000",
    NEXT_TELEMETRY_DISABLED="1",
    DEPLOY_ENV="PRODUCTION",
    EDITION="SELF_HOSTED",
    CONSOLE_API_URL="",
    APP_API_URL="",
    MARKETPLACE_API_URL="http://127.0.0.1:5001",
    MARKETPLACE_URL="http://127.0.0.1:5001"
priority=300
autostart=true
autorestart=unexpected
startsecs=10
stdout_logfile=/data/logs/dify-web.log
stderr_logfile=/data/logs/dify-web.err.log

; ────────────────────────────────────────────────────────────────────────────
; nginx reverse proxy — Dify 웹 게이트웨이 (:18081)
; ────────────────────────────────────────────────────────────────────────────

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=400
autostart=true
autorestart=unexpected
startsecs=3
stdout_logfile=/data/logs/nginx.log
stderr_logfile=/data/logs/nginx.err.log

; ────────────────────────────────────────────────────────────────────────────
; Jenkins 마스터 (:18080 / :50001)
; JENKINS_HOME 은 /data/jenkins (환경변수 리다이렉트).
; 설치 위자드와 CSP 는 docker-compose.override.yaml 의 JAVA_OPTS 와 동일 처리.
; base 이미지의 VOLUME /var/jenkins_home 은 익명 볼륨으로 마운트되지만 사용 안 함.
; ────────────────────────────────────────────────────────────────────────────

[program:jenkins]
; jenkins/jenkins:lts-jdk21 base 의 java 는 /opt/java/openjdk/bin/java ($JAVA_HOME/bin/java).
; /usr/bin/java 는 존재하지 않음.
command=/opt/java/openjdk/bin/java
    -Djenkins.install.runSetupWizard=false
    -Dhudson.model.DirectoryBrowserSupport.CSP=
    -Dhudson.security.csrf.DefaultCrumbIssuer.EXCLUDE_SESSION_ID=true
    -Djenkins.model.Jenkins.slaveAgentPort=50001
    -Djenkins.model.Jenkins.slaveAgentPortEnforce=true
    -jar /usr/share/jenkins/jenkins.war
    --httpPort=18080
    --prefix=
environment=
    JENKINS_HOME="/data/jenkins",
    JAVA_OPTS="-Djenkins.install.runSetupWizard=false"
user=jenkins
priority=500
autostart=true
autorestart=unexpected
startsecs=30
stopwaitsecs=60
stdout_logfile=/data/logs/jenkins.log
stderr_logfile=/data/logs/jenkins.err.log

; [program:jenkins-agent] 제거됨 — Mac 하이브리드 브랜치는 에이전트를 **호스트 Mac**
; 에서 직접 실행하므로 컨테이너 안에서 JNLP agent 를 띄우지 않는다. 호스트 agent
; 연결 가이드: mac-agent-setup.sh + README §2.3. 컨테이너는 JNLP 포트 50001 만
; 노출하고 호스트 agent.jar 가 외부에서 접속한다.
