# e2e-pipeline/offline/supervisord.conf
; ============================================================================
; Zero-Touch QA All-in-One — supervisord 프로세스 매니페스트
; 단일 컨테이너 내부에서 11개 프로세스를 관리한다.
;
; 기동 순서 (priority 낮을수록 먼저):
;   100: 상태 저장소 (postgresql, redis, qdrant, ollama)
;   200: dify-plugin-daemon (postgres 준비 후)
;   300: dify api/worker/worker_beat/web
;   400: nginx (api/web 준비 후)
;   500: jenkins (마스터)
;   600: jenkins-agent (마스터 준비 후 loopback JNLP 연결)
;
; 로그는 /data/logs/<svc>.log 로 분리 저장 (볼륨 지속).
; ============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/data/logs/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ────────────────────────────────────────────────────────────────────────────
; 상태 저장소 (priority 100)
; ────────────────────────────────────────────────────────────────────────────

[program:postgresql]
command=/usr/lib/postgresql/15/bin/postgres -D /data/pg -c config_file=/data/pg/postgresql.conf
user=postgres
priority=100
autostart=true
autorestart=unexpected
stdout_logfile=/data/logs/postgresql.log
stderr_logfile=/data/logs/postgresql.err.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB

[program:redis]
command=/usr/bin/redis-server --dir /data/redis --appendonly yes --bind 127.0.0.1 --port 6379
user=redis
priority=100
autostart=true
autorestart=unexpected
stdout_logfile=/data/logs/redis.log
stderr_logfile=/data/logs/redis.err.log

[program:qdrant]
command=/usr/local/bin/qdrant --storage-dir /data/qdrant
environment=QDRANT__SERVICE__HTTP_PORT="6333",QDRANT__SERVICE__HOST="127.0.0.1"
priority=100
autostart=true
autorestart=unexpected
stdout_logfile=/data/logs/qdrant.log
stderr_logfile=/data/logs/qdrant.err.log

[program:ollama]
command=/usr/bin/ollama serve
environment=OLLAMA_HOST="127.0.0.1:11434",OLLAMA_MODELS="/data/ollama/models"
priority=100
autostart=true
autorestart=unexpected
stdout_logfile=/data/logs/ollama.log
stderr_logfile=/data/logs/ollama.err.log

; ────────────────────────────────────────────────────────────────────────────
; Dify 플러그인 데몬 (priority 200)
; ────────────────────────────────────────────────────────────────────────────

[program:dify-plugin-daemon]
command=/opt/dify-plugin-daemon/plugin_daemon
directory=/opt/dify-plugin-daemon
environment=
    DB_HOST="127.0.0.1",
    DB_PORT="5432",
    DB_USERNAME="postgres",
    DB_PASSWORD="difyai123456",
    DB_DATABASE="dify_plugin",
    REDIS_HOST="127.0.0.1",
    REDIS_PORT="6379",
    SERVER_PORT="5002",
    PLUGIN_WORKING_PATH="/data/dify/plugins/cwd",
    PLUGIN_STORAGE_TYPE="local",
    PLUGIN_STORAGE_LOCAL_ROOT="/data/dify/plugins",
    MAX_PLUGIN_PACKAGE_SIZE="52428800",
    DIFY_INNER_API_URL="http://127.0.0.1:5001",
    DIFY_INNER_API_KEY="QaHbTe77CtuXmsfyhR7+vRjI/+XiGgmVDLZuA+Wf6RTUu8l0mq5kd20MOcc=",
    PPROF_ENABLED="false",
    HOST_WEBHOOK="0.0.0.0"
priority=200
autostart=true
autorestart=unexpected
startsecs=10
stdout_logfile=/data/logs/dify-plugin-daemon.log
stderr_logfile=/data/logs/dify-plugin-daemon.err.log

; ────────────────────────────────────────────────────────────────────────────
; Dify 애플리케이션 서비스 (priority 300)
; Dify 환경변수는 docker-compose.yaml 의 x-shared-api-worker-env 앵커와 동일한
; 조합이되 host 값을 전부 127.0.0.1 로 치환한 것.
; ────────────────────────────────────────────────────────────────────────────

[program:dify-api]
command=/bin/bash /opt/dify-api/docker/entrypoint.sh
directory=/opt/dify-api
environment=
    MODE="api",
    CONSOLE_API_URL="",
    CONSOLE_WEB_URL="",
    SERVICE_API_URL="",
    APP_WEB_URL="",
    DB_HOST="127.0.0.1",
    DB_PORT="5432",
    DB_USERNAME="postgres",
    DB_PASSWORD="difyai123456",
    DB_DATABASE="dify",
    REDIS_HOST="127.0.0.1",
    REDIS_PORT="6379",
    REDIS_PASSWORD="",
    CELERY_BROKER_URL="redis://127.0.0.1:6379/1",
    VECTOR_STORE="qdrant",
    QDRANT_URL="http://127.0.0.1:6333",
    STORAGE_TYPE="opendal",
    OPENDAL_SCHEME="fs",
    OPENDAL_FS_ROOT="/data/dify/storage",
    PLUGIN_DAEMON_URL="http://127.0.0.1:5002",
    PLUGIN_DAEMON_KEY="QaHbTe77CtuXmsfyhR7+vRjI/+XiGgmVDLZuA+Wf6RTUu8l0mq5kd20MOcc=",
    INNER_API_KEY_FOR_PLUGIN="QaHbTe77CtuXmsfyhR7+vRjI/+XiGgmVDLZuA+Wf6RTUu8l0mq5kd20MOcc=",
    MIGRATION_ENABLED="true",
    LOG_LEVEL="INFO",
    SECRET_KEY="dify-allinone-placeholder-CHANGE-ME-via-env",
    SERVER_WORKER_AMOUNT="1",
    SERVER_WORKER_CLASS="gevent"
priority=300
autostart=true
autorestart=unexpected
startsecs=15
stdout_logfile=/data/logs/dify-api.log
stderr_logfile=/data/logs/dify-api.err.log

[program:dify-worker]
command=/bin/bash /opt/dify-api/docker/entrypoint.sh
directory=/opt/dify-api
environment=
    MODE="worker",
    DB_HOST="127.0.0.1",
    DB_PORT="5432",
    DB_USERNAME="postgres",
    DB_PASSWORD="difyai123456",
    DB_DATABASE="dify",
    REDIS_HOST="127.0.0.1",
    REDIS_PORT="6379",
    CELERY_BROKER_URL="redis://127.0.0.1:6379/1",
    VECTOR_STORE="qdrant",
    QDRANT_URL="http://127.0.0.1:6333",
    STORAGE_TYPE="opendal",
    OPENDAL_SCHEME="fs",
    OPENDAL_FS_ROOT="/data/dify/storage",
    PLUGIN_DAEMON_URL="http://127.0.0.1:5002",
    PLUGIN_DAEMON_KEY="QaHbTe77CtuXmsfyhR7+vRjI/+XiGgmVDLZuA+Wf6RTUu8l0mq5kd20MOcc=",
    INNER_API_KEY_FOR_PLUGIN="QaHbTe77CtuXmsfyhR7+vRjI/+XiGgmVDLZuA+Wf6RTUu8l0mq5kd20MOcc=",
    LOG_LEVEL="INFO",
    SECRET_KEY="dify-allinone-placeholder-CHANGE-ME-via-env"
priority=300
autostart=true
autorestart=unexpected
startsecs=15
stdout_logfile=/data/logs/dify-worker.log
stderr_logfile=/data/logs/dify-worker.err.log

[program:dify-worker-beat]
command=/bin/bash /opt/dify-api/docker/entrypoint.sh
directory=/opt/dify-api
environment=
    MODE="beat",
    DB_HOST="127.0.0.1",
    DB_PORT="5432",
    DB_USERNAME="postgres",
    DB_PASSWORD="difyai123456",
    DB_DATABASE="dify",
    REDIS_HOST="127.0.0.1",
    REDIS_PORT="6379",
    CELERY_BROKER_URL="redis://127.0.0.1:6379/1",
    LOG_LEVEL="INFO",
    SECRET_KEY="dify-allinone-placeholder-CHANGE-ME-via-env"
priority=300
autostart=true
autorestart=unexpected
startsecs=15
stdout_logfile=/data/logs/dify-worker-beat.log
stderr_logfile=/data/logs/dify-worker-beat.err.log

[program:dify-web]
command=/bin/sh -c "cd /opt/dify-web && /opt/dify-web/entrypoint.sh"
directory=/opt/dify-web
environment=
    CONSOLE_API_URL="",
    APP_API_URL="",
    HOSTNAME="127.0.0.1",
    PORT="3000"
priority=300
autostart=true
autorestart=unexpected
startsecs=10
stdout_logfile=/data/logs/dify-web.log
stderr_logfile=/data/logs/dify-web.err.log

; ────────────────────────────────────────────────────────────────────────────
; nginx reverse proxy — Dify 웹 게이트웨이 (:18081)
; ────────────────────────────────────────────────────────────────────────────

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=400
autostart=true
autorestart=unexpected
startsecs=3
stdout_logfile=/data/logs/nginx.log
stderr_logfile=/data/logs/nginx.err.log

; ────────────────────────────────────────────────────────────────────────────
; Jenkins 마스터 (:18080 / :50001)
; JENKINS_HOME 은 /data/jenkins (환경변수 리다이렉트).
; 설치 위자드와 CSP 는 docker-compose.override.yaml 의 JAVA_OPTS 와 동일 처리.
; base 이미지의 VOLUME /var/jenkins_home 은 익명 볼륨으로 마운트되지만 사용 안 함.
; ────────────────────────────────────────────────────────────────────────────

[program:jenkins]
command=/usr/bin/java
    -Djenkins.install.runSetupWizard=false
    -Dhudson.model.DirectoryBrowserSupport.CSP=
    -Dhudson.security.csrf.DefaultCrumbIssuer.EXCLUDE_SESSION_ID=true
    -jar /usr/share/jenkins/jenkins.war
    --httpPort=18080
    --prefix=
environment=
    JENKINS_HOME="/data/jenkins",
    JAVA_OPTS="-Djenkins.install.runSetupWizard=false"
user=jenkins
priority=500
autostart=true
autorestart=unexpected
startsecs=30
stopwaitsecs=60
stdout_logfile=/data/logs/jenkins.log
stderr_logfile=/data/logs/jenkins.err.log

; ────────────────────────────────────────────────────────────────────────────
; Jenkins Agent (loopback JNLP) — UI 테스트 headless 실행용
; entrypoint-allinone.sh 가 NODE_SECRET 을 추출해 파일로 떨어뜨린 후 기동한다.
; ────────────────────────────────────────────────────────────────────────────

[program:jenkins-agent]
command=/bin/bash /opt/jenkins-agent-run.sh
environment=
    JENKINS_URL="http://127.0.0.1:18080",
    AGENT_NAME="mac-ui-tester",
    AGENT_WORKDIR="/data/jenkins-agent"
priority=600
autostart=false
autorestart=unexpected
startsecs=10
stdout_logfile=/data/logs/jenkins-agent.log
stderr_logfile=/data/logs/jenkins-agent.err.log
