# scripts/supervisord.allinone.conf
; ============================================================================
; TTC 4-Pipeline All-in-One — supervisord 프로세스 매니페스트
; 단일 컨테이너에서 11개 프로세스 관리 (기존 offline 10개 + SonarQube 추가).
;
; 포트 체계 (기존 Playwright 컨테이너 18080/18081/50001 와 분리):
;   28080 Jenkins HTTP
;   50002 Jenkins agent (JNLP)
;   28081 nginx → Dify gateway
;   29000 SonarQube HTTP
;   5432/6379/6333/5001/5002/3000 은 컨테이너 내부 전용 (호스트 미노출)
;
; 기동 순서:
;   100: postgresql, redis, qdrant                (상태 저장소)
;   150: sonarqube                                (PG 이후, Jenkins 전)
;   200: dify-plugin-daemon
;   300: dify-api / worker / worker_beat / web
;   400: nginx
;   500: jenkins
;
; 로그는 /data/logs/<svc>.log (볼륨 지속).
; ============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/data/logs/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ────────────────────────────────────────────────────────────────────────────
; 상태 저장소 (priority 100)
; ────────────────────────────────────────────────────────────────────────────

[program:postgresql]
command=/usr/lib/postgresql/15/bin/postgres -D /data/pg -c config_file=/data/pg/postgresql.conf
user=postgres
priority=100
autostart=true
autorestart=unexpected
stdout_logfile=/data/logs/postgresql.log
stderr_logfile=/data/logs/postgresql.err.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB

[program:redis]
command=/usr/bin/redis-server --dir /data/redis --appendonly yes --bind 127.0.0.1 --port 6379
user=redis
priority=100
autostart=true
autorestart=unexpected
stdout_logfile=/data/logs/redis.log
stderr_logfile=/data/logs/redis.err.log

[program:qdrant]
command=/usr/local/bin/qdrant
environment=
    QDRANT__SERVICE__HTTP_PORT="6333",
    QDRANT__SERVICE__HOST="127.0.0.1",
    QDRANT__STORAGE__STORAGE_PATH="/data/qdrant",
    QDRANT__STORAGE__SNAPSHOTS_PATH="/data/qdrant/snapshots",
    QDRANT__TELEMETRY_DISABLED="true"
priority=100
autostart=true
autorestart=unexpected
stdout_logfile=/data/logs/qdrant.log
stderr_logfile=/data/logs/qdrant.err.log

; ────────────────────────────────────────────────────────────────────────────
; SonarQube (priority 150) — 파이프라인 2 정적분석 엔진
; 내장 Elasticsearch 가 vm.max_map_count >= 262144 을 기본 요구하나,
; sonar.es.bootstrap.checks.disable=true (Dockerfile 에서 삽입) 로 우회.
; PG 는 127.0.0.1:5432/sonar 를 사용 (pg-init 에서 생성).
; ────────────────────────────────────────────────────────────────────────────

[program:sonarqube]
command=/opt/sonarqube/bin/sonar.sh console
directory=/opt/sonarqube
user=sonar
priority=150
autostart=false
autorestart=unexpected
startsecs=60
stopwaitsecs=60
environment=
    SONAR_JDBC_URL="jdbc:postgresql://127.0.0.1:5432/sonar",
    SONAR_JDBC_USERNAME="sonar",
    SONAR_JDBC_PASSWORD="sonar",
    SONAR_WEB_PORT="9000",
    SONAR_WEB_HOST="0.0.0.0",
    SONAR_SEARCH_JAVAADDITIONALOPTS="-Dnode.store.allow_mmap=false"
stdout_logfile=/data/logs/sonarqube.log
stderr_logfile=/data/logs/sonarqube.err.log

; sonar.sh console 은 기본 :9000 에 바인딩. 호스트 29000 매핑은 docker run/compose 에서.

; ────────────────────────────────────────────────────────────────────────────
; Dify plugin daemon (priority 200)
; ────────────────────────────────────────────────────────────────────────────

[program:dify-plugin-daemon]
command=/opt/dify-plugin-daemon/main
directory=/opt/dify-plugin-daemon
environment=
    PATH="/usr/local/bin:/usr/bin:/bin",
    PLATFORM="local",
    SERVER_KEY="QaHbTe77CtuXmsfyhR7+vRjI/+XiGgmVDLZuA+Wf6RTUu8l0mq5kd20MOcc=",
    UV_PATH="/usr/local/bin/uv",
    TIKTOKEN_CACHE_DIR="/opt/dify-plugin-daemon/.tiktoken",
    DB_HOST="127.0.0.1",
    DB_PORT="5432",
    DB_USERNAME="postgres",
    DB_PASSWORD="difyai123456",
    DB_DATABASE="dify_plugin",
    REDIS_HOST="127.0.0.1",
    REDIS_PORT="6379",
    SERVER_PORT="5002",
    PLUGIN_WORKING_PATH="/data/dify/plugins/cwd",
    PLUGIN_STORAGE_TYPE="local",
    PLUGIN_STORAGE_LOCAL_ROOT="/data/dify/plugins",
    PLUGIN_INSTALLED_PATH="plugin",
    PLUGIN_PACKAGE_CACHE_PATH="plugin_packages",
    PLUGIN_MEDIA_CACHE_PATH="assets",
    PLUGIN_REMOTE_INSTALLING_HOST="0.0.0.0",
    PLUGIN_REMOTE_INSTALLING_PORT="5003",
    MAX_PLUGIN_PACKAGE_SIZE="52428800",
    PYTHON_ENV_INIT_TIMEOUT="120",
    PLUGIN_MAX_EXECUTION_TIMEOUT="600",
    FORCE_VERIFYING_SIGNATURE="false",
    DIFY_INNER_API_URL="http://127.0.0.1:5001",
    DIFY_INNER_API_KEY="QaHbTe77CtuXmsfyhR7+vRjI/+XiGgmVDLZuA+Wf6RTUu8l0mq5kd20MOcc=",
    PPROF_ENABLED="false",
    HOST_WEBHOOK="0.0.0.0"
priority=200
autostart=true
autorestart=unexpected
startsecs=10
stdout_logfile=/data/logs/dify-plugin-daemon.log
stderr_logfile=/data/logs/dify-plugin-daemon.err.log

; ────────────────────────────────────────────────────────────────────────────
; Dify 애플리케이션 (priority 300)
; ────────────────────────────────────────────────────────────────────────────

[program:dify-api]
command=/bin/bash /opt/dify-api/api/docker/entrypoint.sh
directory=/opt/dify-api/api
environment=
    PATH="/opt/dify-api/api/.venv/bin:/usr/local/bin:/usr/bin:/bin",
    VIRTUAL_ENV="/opt/dify-api/api/.venv",
    FLASK_APP="app.py",
    TIKTOKEN_CACHE_DIR="/opt/dify-api/api/.tiktoken_cache",
    NLTK_DATA="/usr/local/share/nltk_data",
    LANG="C.UTF-8",
    LC_ALL="C.UTF-8",
    PYTHONIOENCODING="utf-8",
    EDITION="SELF_HOSTED",
    DEPLOY_ENV="PRODUCTION",
    MODE="api",
    CONSOLE_API_URL="",
    CONSOLE_WEB_URL="",
    SERVICE_API_URL="",
    APP_WEB_URL="",
    DB_HOST="127.0.0.1",
    DB_PORT="5432",
    DB_USERNAME="postgres",
    DB_PASSWORD="difyai123456",
    DB_DATABASE="dify",
    REDIS_HOST="127.0.0.1",
    REDIS_PORT="6379",
    REDIS_PASSWORD="",
    CELERY_BROKER_URL="redis://127.0.0.1:6379/1",
    VECTOR_STORE="qdrant",
    QDRANT_URL="http://127.0.0.1:6333",
    STORAGE_TYPE="opendal",
    OPENDAL_SCHEME="fs",
    OPENDAL_FS_ROOT="/data/dify/storage",
    PLUGIN_DAEMON_URL="http://127.0.0.1:5002",
    PLUGIN_DAEMON_KEY="QaHbTe77CtuXmsfyhR7+vRjI/+XiGgmVDLZuA+Wf6RTUu8l0mq5kd20MOcc=",
    INNER_API_KEY_FOR_PLUGIN="QaHbTe77CtuXmsfyhR7+vRjI/+XiGgmVDLZuA+Wf6RTUu8l0mq5kd20MOcc=",
    MIGRATION_ENABLED="true",
    LOG_LEVEL="INFO",
    SECRET_KEY="dify-allinone-placeholder-CHANGE-ME-via-env",
    SERVER_WORKER_AMOUNT="1",
    SERVER_WORKER_CLASS="gevent",
    SERVER_WORKER_CONNECTIONS="1000",
    GUNICORN_TIMEOUT="360"
priority=300
; entrypoint 가 plugin-daemon 헬스체크 후 명시적으로 start
autostart=false
autorestart=unexpected
startsecs=15
stdout_logfile=/data/logs/dify-api.log
stderr_logfile=/data/logs/dify-api.err.log

[program:dify-worker]
command=/bin/bash /opt/dify-api/api/docker/entrypoint.sh
directory=/opt/dify-api/api
environment=
    PATH="/opt/dify-api/api/.venv/bin:/usr/local/bin:/usr/bin:/bin",
    VIRTUAL_ENV="/opt/dify-api/api/.venv",
    FLASK_APP="app.py",
    TIKTOKEN_CACHE_DIR="/opt/dify-api/api/.tiktoken_cache",
    NLTK_DATA="/usr/local/share/nltk_data",
    LANG="C.UTF-8",
    LC_ALL="C.UTF-8",
    PYTHONIOENCODING="utf-8",
    EDITION="SELF_HOSTED",
    DEPLOY_ENV="PRODUCTION",
    MODE="worker",
    DB_HOST="127.0.0.1",
    DB_PORT="5432",
    DB_USERNAME="postgres",
    DB_PASSWORD="difyai123456",
    DB_DATABASE="dify",
    REDIS_HOST="127.0.0.1",
    REDIS_PORT="6379",
    CELERY_BROKER_URL="redis://127.0.0.1:6379/1",
    VECTOR_STORE="qdrant",
    QDRANT_URL="http://127.0.0.1:6333",
    STORAGE_TYPE="opendal",
    OPENDAL_SCHEME="fs",
    OPENDAL_FS_ROOT="/data/dify/storage",
    PLUGIN_DAEMON_URL="http://127.0.0.1:5002",
    PLUGIN_DAEMON_KEY="QaHbTe77CtuXmsfyhR7+vRjI/+XiGgmVDLZuA+Wf6RTUu8l0mq5kd20MOcc=",
    INNER_API_KEY_FOR_PLUGIN="QaHbTe77CtuXmsfyhR7+vRjI/+XiGgmVDLZuA+Wf6RTUu8l0mq5kd20MOcc=",
    LOG_LEVEL="INFO",
    SECRET_KEY="dify-allinone-placeholder-CHANGE-ME-via-env"
priority=300
autostart=true
autorestart=unexpected
startsecs=15
stdout_logfile=/data/logs/dify-worker.log
stderr_logfile=/data/logs/dify-worker.err.log

[program:dify-worker-beat]
command=/bin/bash /opt/dify-api/api/docker/entrypoint.sh
directory=/opt/dify-api/api
environment=
    PATH="/opt/dify-api/api/.venv/bin:/usr/local/bin:/usr/bin:/bin",
    VIRTUAL_ENV="/opt/dify-api/api/.venv",
    FLASK_APP="app.py",
    TIKTOKEN_CACHE_DIR="/opt/dify-api/api/.tiktoken_cache",
    NLTK_DATA="/usr/local/share/nltk_data",
    LANG="C.UTF-8",
    LC_ALL="C.UTF-8",
    PYTHONIOENCODING="utf-8",
    EDITION="SELF_HOSTED",
    DEPLOY_ENV="PRODUCTION",
    MODE="beat",
    DB_HOST="127.0.0.1",
    DB_PORT="5432",
    DB_USERNAME="postgres",
    DB_PASSWORD="difyai123456",
    DB_DATABASE="dify",
    REDIS_HOST="127.0.0.1",
    REDIS_PORT="6379",
    CELERY_BROKER_URL="redis://127.0.0.1:6379/1",
    LOG_LEVEL="INFO",
    SECRET_KEY="dify-allinone-placeholder-CHANGE-ME-via-env"
priority=300
autostart=true
autorestart=unexpected
startsecs=15
stdout_logfile=/data/logs/dify-worker-beat.log
stderr_logfile=/data/logs/dify-worker-beat.err.log

[program:dify-web]
command=/bin/sh -c "cd /opt/dify-web/web && /opt/dify-web/web/entrypoint.sh"
directory=/opt/dify-web/web
environment=
    PATH="/usr/local/bin:/usr/bin:/bin",
    NODE_ENV="production",
    HOSTNAME="127.0.0.1",
    PORT="3000",
    NEXT_TELEMETRY_DISABLED="1",
    DEPLOY_ENV="PRODUCTION",
    EDITION="SELF_HOSTED",
    CONSOLE_API_URL="",
    APP_API_URL="",
    MARKETPLACE_API_URL="http://127.0.0.1:5001",
    MARKETPLACE_URL="http://127.0.0.1:5001"
priority=300
autostart=true
autorestart=unexpected
startsecs=10
stdout_logfile=/data/logs/dify-web.log
stderr_logfile=/data/logs/dify-web.err.log

; ────────────────────────────────────────────────────────────────────────────
; nginx reverse proxy — Dify 웹 게이트웨이 (:28081)
; ────────────────────────────────────────────────────────────────────────────

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=400
autostart=true
autorestart=unexpected
startsecs=3
stdout_logfile=/data/logs/nginx.log
stderr_logfile=/data/logs/nginx.err.log

; ────────────────────────────────────────────────────────────────────────────
; Jenkins 마스터 (:28080 / :50002)
; ────────────────────────────────────────────────────────────────────────────

[program:jenkins]
command=/opt/java/openjdk/bin/java
    -Djenkins.home=/data/jenkins
    -Djenkins.install.runSetupWizard=false
    -Dhudson.model.DirectoryBrowserSupport.CSP=
    -Dhudson.security.csrf.DefaultCrumbIssuer.EXCLUDE_SESSION_ID=true
    -Djenkins.model.Jenkins.slaveAgentPort=50002
    -Djenkins.model.Jenkins.slaveAgentPortEnforce=true
    -jar /usr/share/jenkins/jenkins.war
    --httpPort=28080
    --prefix=
environment=
    JENKINS_HOME="/data/jenkins",
    JAVA_OPTS="-Djenkins.install.runSetupWizard=false"
user=jenkins
priority=500
autostart=true
autorestart=unexpected
startsecs=30
stopwaitsecs=60
stdout_logfile=/data/logs/jenkins.log
stderr_logfile=/data/logs/jenkins.err.log
