# e2e-pipeline/perf-offline/supervisord.conf
; ============================================================================
; JMeter Performance All-in-One — supervisord 프로세스 매니페스트
;
; Ollama 는 호스트에서 동작 (host.docker.internal:11434) — 컨테이너 안 미실행.
;
; 기동 순서 (priority 낮을수록 먼저):
;   200: xvfb (가상 디스플레이)
;   210: fluxbox (경량 윈도우 매니저)
;   220: jmeter-gui (Apache JMeter Swing GUI) — JMETER_GUI_AUTOSTART 환경변수로 토글
;   230: x11vnc (가상 디스플레이 → VNC 서버)
;   240: novnc (VNC → WebSocket → 브라우저 노출)
;   300: nginx (단일 HTTP 게이트웨이 :18090)
;
; 로그는 /data/logs/<svc>.log 로 분리 저장 (볼륨 지속).
; ============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/data/logs/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ────────────────────────────────────────────────────────────────────────────
; 가상 디스플레이 스택
; ────────────────────────────────────────────────────────────────────────────

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1600x900x24 -nolisten tcp -dpi 96
priority=200
autostart=true
autorestart=unexpected
stdout_logfile=/data/logs/xvfb.log
stderr_logfile=/data/logs/xvfb.err.log

[program:fluxbox]
command=/usr/bin/fluxbox -display :1
environment=DISPLAY=":1",HOME="/data/jmeter"
user=root
priority=210
autostart=true
autorestart=unexpected
startsecs=2
stdout_logfile=/data/logs/fluxbox.log
stderr_logfile=/data/logs/fluxbox.err.log

; ────────────────────────────────────────────────────────────────────────────
; JMeter GUI — 평상시 시나리오 작성/디버깅용
; entrypoint-allinone.sh 가 JMETER_GUI_AUTOSTART 환경변수에 따라 supervisorctl 로
; 시작/중지하므로 여기서는 autostart=false 로 두고 entrypoint 가 명시적으로 제어.
; ────────────────────────────────────────────────────────────────────────────

[program:jmeter-gui]
command=/opt/apache-jmeter/bin/jmeter -p /opt/apache-jmeter/bin/user.properties
environment=
    DISPLAY=":1",
    HOME="/data/jmeter",
    JMETER_HOME="/opt/apache-jmeter",
    JVM_ARGS="-Xms512m -Xmx2g"
directory=/data/jmeter
user=root
priority=220
autostart=false
autorestart=true
startsecs=10
stdout_logfile=/data/logs/jmeter-gui.log
stderr_logfile=/data/logs/jmeter-gui.err.log

; ────────────────────────────────────────────────────────────────────────────
; VNC + noVNC — 가상 디스플레이를 브라우저로 노출
; ────────────────────────────────────────────────────────────────────────────

[program:x11vnc]
; -nopw : 패스워드 없음 (사내망 전용; 외부 노출 시 -passwd 또는 -rfbauth 권장)
; -shared : 다중 클라이언트 동시 접속 허용
; -forever : 클라이언트 연결 끊겨도 종료 안 함
command=/usr/bin/x11vnc -display :1 -forever -shared -nopw -rfbport 5900 -quiet
priority=230
autostart=true
autorestart=unexpected
startsecs=3
stdout_logfile=/data/logs/x11vnc.log
stderr_logfile=/data/logs/x11vnc.err.log

[program:novnc]
command=/usr/local/bin/websockify --web=/opt/novnc 6080 127.0.0.1:5900
priority=240
autostart=true
autorestart=unexpected
startsecs=3
stdout_logfile=/data/logs/novnc.log
stderr_logfile=/data/logs/novnc.err.log

; ────────────────────────────────────────────────────────────────────────────
; nginx — 통합 HTTP 게이트웨이 (:18090)
; / → noVNC, /ollama/ → Ollama API
; ────────────────────────────────────────────────────────────────────────────

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=300
autostart=true
autorestart=unexpected
startsecs=3
stdout_logfile=/data/logs/nginx.log
stderr_logfile=/data/logs/nginx.err.log
